# touch-keyboard-handler.service
[Unit]
Description=Touch keyboard handler

[Service]
WorkingDirectory=/etc/touch_keyboard
ExecStart=/usr/sbin/touch_keyboard_handler
Type=simple

